# patches/dropbear.conffiles
--- a/debian/dropbear.conffiles	2025-05-07 16:30:00.000000000 +0200
+++ b/debian/dropbear.conffiles	2025-07-17 15:13:33.331994616 +0200
@@ -1,3 +1 @@
-/etc/init.d/dropbear
-/etc/dropbear/run
-/etc/dropbear/log/run
+/etc/init.d/dropbear-static
